# response_differences/service/named/run
#!/bin/bash
exec 2>&1
exec /usr/sbin/named -u "$(id -u nobody)" -g -c named.conf
